# oasis-dist/src/deb/control/postinst
#! /bin/sh

update-rc.d [[name]] defaults > /dev/null

invoke-rc.d [[name]] start
