# oasis-dist/src/deb/control/postrm
#! /bin/sh

update-rc.d @name@ remove > /dev/null

#TODO:
# remove user and group on purge
